# openwrt/xiaomi_mi-router-3/files/etc/uci-defaults/50-password.sh
#!/usr/bin/env sh

set -ex

passwd=""
user="admin"

printf "%s\n%s" "${passwd}" "${passwd}" | passwd "${user}"
