# home/private_dot_config/autostart/50-chatterino.desktop
[Desktop Entry]
Exec=chatterino
Name=chatterino
Type=Application
