# home/private_dot_config/autostart/50-discord.desktop
[Desktop Entry]
Exec=discord --ozone-platform-hint=wayland
Name=discord
Type=Application
