# home/private_dot_config/autostart/50-keyring.desktop
[Desktop Entry]
Exec=gnome-keyring-daemon --login --components=secrets
Name=keyring
Type=Application
