# home/private_dot_config/autostart/50-telegram.desktop
[Desktop Entry]
Exec=telegram-desktop
Name=telegram
Type=Application
