# root/etc/xdg/autostart/40-dbus.desktop
[Desktop Entry]
Exec=dbus-update-activation-environment --all --systemd
Name=dbus
Type=Application
