# root/etc/xdg/autostart/50-clipboard-image.desktop
[Desktop Entry]
Exec=wl-paste --type image --watch cliphist store
Name=clipboard-image
Type=Application
