# root/etc/xdg/autostart/50-clipboard.desktop
[Desktop Entry]
Exec=wl-paste --type image --watch cliphist store && wl-paste --type text --watch cliphist store
Name=clipboard
Type=Application
