# root/etc/xdg/autostart/50-clipboard-text.desktop
[Desktop Entry]
Exec=wl-paste --type text --watch cliphist store
Name=clipboard-text
Type=Application
