# root/etc/xdg/autostart/50-keyring.desktop
[Desktop Entry]
Exec=gnome-keyring-daemon --start --components=secrets
Name=keyring
Type=Application
