# root/etc/xdg/autostart/50-kdeconnect.desktop
[Desktop Entry]
Exec=/usr/lib/kdeconnectd
Name=kdeconnect
Type=Application
