# root/etc/xdg/autostart/50-network.desktop
[Desktop Entry]
Exec=nm-applet
Name=network
Type=Application
