# root/etc/xdg/autostart/50-playerctl.desktop
[Desktop Entry]
Exec=playerctld daemon
Name=playerctl
Type=Application
