# root/etc/xdg/autostart/50-playerctl.desktop
[Desktop Entry]
Exec=playerctld
Name=playerctl
Type=Application
