# root/etc/xdg/autostart/50-swww-daemon.desktop
[Desktop Entry]
Exec=swww-daemon
Name=swww-daemon
Type=Application
