# root/etc/xdg/autostart/50-swww.desktop
[Desktop Entry]
Exec=swww-daemon && swww.sh ${XDG_PICTURES_DIR}/Wallpapers
Name=swww
Type=Application
