# openwrt/tplink_archer-c6-v2/files/etc/uci-defaults/50-password.sh
#!/usr/bin/env sh

set -ex

PASSWD=
USER=root

printf "%s\n%s" "${PASSWD}" "${PASSWD}" | passwd "${USER}"
