# openwrt/tplink_archer-c6-v2/files/etc/uci-defaults/50-password.sh
#!/usr/bin/env sh

set -ex

passwd=
user=root

printf "%s\n%s" "${passwd}" "${passwd}" | passwd "${user}"
